# mkosi.extra/usr/lib/udev/rules.d/60-usb-dirty-pages-udev.rules
ACTION=="add", KERNEL=="sd[a-z]", SUBSYSTEM=="block", ENV{ID_USB_TYPE}=="disk", RUN+="/usr/lib/usb-dirty-pages-udev %k"
