# mkosi.extra/usr/lib/systemd/system/plasma-setup.service.d/kdelinux-firstboot.conf
# SPDX-License-Identifier: GPL-2.0-only OR GPL-3.0-only OR LicenseRef-KDE-Accepted-GPL
#
# SPDX-FileCopyrightText: 2026 Kristen McWilliam <kristen@kde.org>

# Only run plasma-setup when the installer marker exists, to prevent users with
# existing installations from being affected.
[Unit]
ConditionPathExists=/var/lib/kde-linux/installed-with-plasma-setup
