# mkosi.extra/usr/lib/systemd/system/kde-linux-install-new-flatpak-apps.service
# SPDX-License-Identifier: GPL-2.0-only OR GPL-3.0-only OR LicenseRef-KDE-Accepted-GPL
# SPDX-FileCopyrightText: 2026 Nate Graham <nate@kde.org>

[Unit]
Description=Install any new Flatpak apps added to booted OS version
ConditionPathIsReadWrite=/etc
ConditionNeedsUpdate=|/etc
DefaultDependencies=no
After=local-fs.target
Before=sysinit.target systemd-update-done.service
Conflicts=shutdown.target
Before=shutdown.target

[Service]
Type=oneshot
ExecStart=/usr/lib/install-new-flatpaks
