# mkosi.extra/usr/lib/systemd/system/sysext-trigger.service
# SPDX-License-Identifier: GPL-2.0-only OR GPL-3.0-only OR LicenseRef-KDE-Accepted-GPL
# SPDX-FileCopyrightText: 2026 Hadi Chokr <hadichokr@icloud.com>

[Unit]
Description=Prepare sysext trigger and routing (installed system only)
DefaultDependencies=no
After=local-fs.target
Before=systemd-sysext.service
ConditionPathExists=!/live/
ConditionKernelCommandLine=!kde-linux.live=1

[Service]
Type=oneshot

# --- Ensure base dirs exist ---
ExecStart=/usr/bin/mkdir -p /var/lib/extensions
ExecStart=/usr/bin/mkdir -p /var/lib/extensions.mutable/usr

# --- Ensure correct permissions ---
ExecStart=/usr/bin/chmod 0755 /var/lib/extensions.mutable
ExecStart=/usr/bin/chmod 0755 /var/lib/extensions.mutable/usr
ExecStart=/usr/bin/chmod 0755 /var/lib/extensions.mutable/opt

# --- Route /usr/local → /opt/local ---
ExecStart=/usr/bin/ln -sfn /opt/local /var/lib/extensions.mutable/usr/local

# --- Trigger sysext normally ---
ExecStart=/usr/bin/touch /var/lib/extensions/auto-enable

[Install]
WantedBy=sysinit.target
