# mkosi.extra/usr/lib/systemd/system/systemd-suspend.service.d/10-nvidia-no-freeze-session.conf
# SPDX-License-Identifier: CC0-1.0
# SPDX-FileCopyrightText: none

# This shadows the file of the same name provided by NVIDIA and instead reads the below environment
# file every time suspend is attempted. This is so SYSTEMD_SLEEP_FREEZE_USER_SESSIONS is only set to
# false when the nvidia kernel module is loaded, unlike the blanket rule shipped by NVIDIA.

[Service]
EnvironmentFile=-/run/nvidia-no-freeze-session
