# etc/systemd/system/btrfs-snapshot-prune.service
[Unit]
Description=Prune old Btrfs snapshots

[Service]
Type=oneshot
ExecStart=/usr/local/bin/btrfs-snapshot-prune.sh
